# spotify/rootfs/etc/services.d/spotifyd/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Spotify Connect
# Runs librespot
# ==============================================================================
declare -a options
declare name
declare bitrate
declare initial_volume
declare volume_ctrl
declare volume_normalization
declare normalization_method
declare autoplay
declare disable_discovery
declare password
declare username
declare autoplay
declare initial_volume

bashio::log.info 'Starting the Spotify daemon...'

# Bitrate
options+=(--bitrate $(bashio::config 'bitrate'))

# Device name
name=$(bashio::config 'name')
options+=(--name "${name}")

# Bitrate
options+=(--bitrate $(bashio::config 'bitrate'))

# Initial volume
options+=(--initial-volume $(bashio::config 'initial_volume'))

# Volume Control Type
options+=(--volume-ctrl $(bashio::config 'volume_ctrl'))

# Volume Normalization
if bashio::config.true 'volume_normalization'; then
  options+=(--enable-volume-normalisation)
fi

# Normalization Method
options+=(--normalisation-method $(bashio::config 'normalization_method'))

# Autoplay
if bashio::config.true 'autoplay'; then
  options+=(--autoplay)
fi

# Disable Discovery
if bashio::config.true 'disable_discovery'; then
  options+=(--disable-discovery)
fi

# Save writes
options+=(--disable-audio-cache)

# Run PROGRAM set by --onevent before sink is opened and after it is closed.
options+=(--onevent 'sh /usr/local/bin/event.sh')
options+=(--emit-sink-events)

# Are we running in debug mode?
if bashio::debug; then
  options+=(--verbose)
  bashio::log.info "options w/o username and password"
  bashio::log.info "${options[@]}"
fi

# Username / password
if bashio::config.has_value 'username'; then
  bashio::config.require.password

  username=$(bashio::config 'username')
  password=$(bashio::config 'password')

  options+=(--username "${username}")
  options+=(--password "${password}")
fi

# Run librespot
exec librespot "${options[@]}"
